# systemd/lostark-item-checker.service
[Unit]
Description=Lost Ark Item Checker Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=ahrrri
WorkingDirectory=/home/ahrrri/find-angel
Environment=PATH=/home/ahrrri/find-angel/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/ahrrri/find-angel/venv/bin/python -u -m src.core.async_item_checker
Restart=always
RestartSec=10
StandardOutput=append:/home/ahrrri/find-angel/logs/item-checker.log
StandardError=append:/home/ahrrri/find-angel/logs/item-checker-error.log

# 우선순위 및 리소스 제한 (실시간성 중요)
Nice=2
IOSchedulingClass=2
IOSchedulingPriority=2
MemoryMax=1G
CPUQuota=60%

[Install]
WantedBy=multi-user.target